# docker/crontab
* 0 */6 * * *     /usr/local/bin/crawl > /var/log/crawl.log
