# docker/crontab
PATH=/usr/local/bin:/usr/bin:/bin:/app/bin

10 */2 * * *     collect
0  4   * * *     find /var/log/feeder -name '*.log' -type f -mtime +7 -delete
