# docker/runit/log/run
#!/bin/bash

exec 2>&1
exec svlogd -tt "$AARON_LOG_DIR"
